# modules/profile/templates/transmission/transmission-daemon.service.erb
[Unit]
Description=Transmission BitTorrent Daemon
After=network.target

[Service]
User=debian-transmission
Type=notify
Environment=TRANSMISSION_WEB_HOME=<%= @web_directory %>
ExecStart=/usr/bin/transmission-daemon -f --log-error
ExecReload=/bin/kill -s HUP $MAINPID

[Install]
WantedBy=multi-user.target
